# .config/systemd/user/himalaya-watcher.service
[Unit]
Description=Himalaya account watcher

[Service]
ExecStart=/usr/bin/himalaya watch
Environment=RUST_LOG=debug
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
